# deploy/fintrax-queue.service
# Fintrax queue worker — install to /etc/systemd/system/fintrax-queue.service
# then: systemctl daemon-reload && systemctl enable --now fintrax-queue
# The deploy workflow runs `systemctl restart fintrax-queue` so the worker
# picks up the new code after each release flip.
[Unit]
Description=Fintrax queue worker
After=network.target

[Service]
User=www-data
Group=www-data
Restart=always
RestartSec=5
WorkingDirectory=/opt/fintrax/current
ExecStart=/usr/bin/php8.5 /opt/fintrax/current/artisan queue:work --sleep=3 --tries=3 --max-time=3600

[Install]
WantedBy=multi-user.target
